# Asawa/main_service_restart.service
[Unit]
Description=Restart the run_main

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart run_main.service