# Devesh Electricals/edge_device/edge_device/boot_setup.timer
[Unit]
Description=hotspot timer
[Timer]
OnBootSec=1min
[Install]
WantedBy=timers.target